# huginn/systemd/huginn.service
[Unit]
Description=Huginn v2 AI Daemon
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/.local/share/huginn
Environment="PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=uv run --project %h/dotfiles/huginn python3 %h/dotfiles/huginn/v2/daemon.py
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
