# huginn/systemd/huginn-morning.timer
[Unit]
Description=Huginn morning briefing timer
Requires=huginn.service

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true

[Install]
WantedBy=timers.target
